# var/lib/opkg/info/autofs.list
/usr/lib/autofs/mount_generic.so	0100755
/usr/sbin/automount	0100755
/usr/lib/autofs/mount_afs.so	0100755
/usr/lib/autofs/mount_bind.so	0100755
/usr/lib/autofs/parse_sun.so	0100755
/usr/lib/autofs/mount_nfs.so	0100755
/etc/auto.net	0100755
/usr/lib/autofs/lookup_files.so	0120777
/usr/lib/autofs/mount_autofs.so	0100755
/usr/lib/autofs/parse_amd.so	0100755
/usr/lib/autofs/lookup_dir.so	0100755
/etc/auto.master	0100644
/usr/lib/autofs/lookup_file.so	0100755
/usr/lib/autofs/lookup_yp.so	0100755
/usr/lib/libautofs.so	0100755
/etc/init.d/autofs	0100755
/usr/lib/autofs/mount_changer.so	0100755
/usr/lib/autofs/lookup_program.so	0100755
/usr/lib/autofs/lookup_hosts.so	0100755
/usr/lib/autofs	040755
/etc/default/autofs	0100644
/etc/autofs.conf	0100644
/etc/auto.network	0100644
/etc/default/volatiles/99_autofs	0100644
/usr/lib/autofs/lookup_multi.so	0100755
/usr/lib/autofs/mount_nfs4.so	0120777
/usr/lib/autofs/lookup_userhome.so	0100755
/usr/lib/autofs/lookup_nis.so	0120777
