# var/lib/opkg/info/busybox-telnetd.list
/etc/init.d	040755
/etc	040755
/etc/init.d/telnetd.busybox	0100755
