# var/lib/opkg/info/cronie.conffiles
/etc/crontab
/etc/cron.deny
/etc/crontab
/etc/default/crond
