# var/lib/opkg/info/cronie.list
/var	040755
/etc/cron.weekly	040755
/etc/init.d	040755
/usr/sbin	040755
/etc/cron.deny	0100644
/etc/cron	040770
/usr/sbin/crond	0100755
/etc/crontab	0100600
/etc/init.d/crond	0100755
/usr/bin	040755
/usr	040755
/etc	040755
/etc/default	040755
/etc/cron.monthly	040755
/var/spool	040755
/etc/pam.d	040755
/etc/cron.daily	040755
/etc/default/crond	0100644
/usr/bin/cronnext	0100755
/etc/cron/crontabs	040770
/usr/bin/crontab	0102755
/etc/pam.d/crond	0100644
/etc/cron.hourly	040755
/etc/cron.d	040755
