# var/lib/opkg/info/fake-hwclock.postinst
#!/bin/sh
set -e
if [ -n "$D" ]; then
        [[ -f $D/etc/fake-hwclock.data ]] || date -u '+%Y-%m-%d %H:%M:%S' > $D/etc/fake-hwclock.data
fi
if true && type update-rc.d >/dev/null 2>/dev/null; then
	if [ -n "$D" ]; then
		OPT="-r $D"
	else
		OPT="-s"
	fi
	update-rc.d $OPT fake-hwclock start 01 S . stop 25 0 6 .
fi
