# var/lib/opkg/info/gigablue-libreader-gbmv200.list
/etc/init.d	040755
/etc/init.d/gigablue-libreader	0100755
/etc	040755
