# var/lib/opkg/info/gigablue-loadmodules-gbmv200.list
/etc/init.d	040755
/etc	040755
/etc/init.d/gigablue-loadmodules	0100755
