# var/lib/opkg/info/initscripts.list
/etc/rc3.d/S15mountnfs.sh	0120777
/etc/default/devpts	0100755
/etc/rc4.d/S50softcam	0120777
/etc/rcS.d/S03mountall.sh	0120777
/etc/init.d/sysfs.sh	0100755
/etc/init.d/umountfs	0100755
/etc/rc1.d/K31umountnfs.sh	0120777
/etc/rcS.d/S38urandom	0120777
/etc/rc5.d/S15mountnfs.sh	0120777
/etc/init.d/read-only-rootfs-hook.sh	0100755
/etc/init.d/checkroot.sh	0100755
/etc/rc0.d/K50softcam	0120777
/etc/rc0.d/S20sendsigs	0120777
/etc/rc0.d/S40umountfs	0120777
/etc/rc2.d/S99rmnologin.sh	0120777
/etc/default/volatiles/00_core	0100644
/var/lib/urandom	040755
/etc/rc2.d/S50softcam	0120777
/etc/init.d/banner.sh	0100755
/etc/rcS.d/S05checkroot.sh	0120777
/etc/rc6.d/S38urandom	0120777
/etc/rcS.d/S29read-only-rootfs-hook.sh	0120777
/etc/init.d/populate-volatile.sh	0100755
/etc/init.d/softcam.None	0100755
/etc/rc4.d/S99rmnologin.sh	0120777
/etc/init.d/bootmisc.sh	0100755
/etc/rcS.d/S06alignment.sh	0120777
/etc/rcS.d/S06hotplug.sh	0120777
/etc/rcS.d/S02banner.sh	0120777
/etc/init.d/devpts.sh	0100755
/etc/init.d/urandom	0100755
/etc/rc0.d/S90halt	0120777
/etc/rc5.d/S50softcam	0120777
/etc/rc2.d/S15mountnfs.sh	0120777
/etc/rc0.d/K31umountnfs.sh	0120777
/etc/rc0.d/S25save-rtc.sh	0120777
/etc/init.d/rmnologin.sh	0100755
/etc/init.d/halt	0100755
/etc/rc1.d/K50softcam	0120777
/etc/rc4.d/S15mountnfs.sh	0120777
/etc/init.d/hostname.sh	0100755
/etc/rc6.d/K50softcam	0120777
/etc/init.d/sendsigs	0100755
/etc/rcS.d/S37populate-volatile.sh	0120777
/etc/rc3.d/S50softcam	0120777
/etc/init.d/save-rtc.sh	0100755
/etc/init.d/hotplug.sh	0100755
/etc/rcS.d/S06devpts.sh	0120777
/etc/rc6.d/S25save-rtc.sh	0120777
/etc/rc3.d/S99rmnologin.sh	0120777
/etc/rcS.d/S95cardserver	0120777
/etc/logrotate-dmesg.conf	0100644
/etc/init.d/dmesg.sh	0100755
/etc/rc6.d/S20sendsigs	0120777
/etc/rc6.d/S40umountfs	0120777
/etc/rcS.d/S39hostname.sh	0120777
/etc/rc6.d/S90reboot	0120777
/etc/rc0.d/S38urandom	0120777
/etc/init.d/reboot	0100755
/etc/init.d/alignment.sh	0100755
/etc/rc5.d/S99rmnologin.sh	0120777
/etc/rc6.d/K31umountnfs.sh	0120777
/etc/rcS.d/S02sysfs.sh	0120777
/etc/init.d/umountnfs.sh	0100755
/etc/rcS.d/S55bootmisc.sh	0120777
/etc/init.d/mountnfs.sh	0100755
/etc/init.d/cardserver.None	0100755
/etc/rcS.d/S38dmesg.sh	0120777
/etc/init.d/mountall.sh	0100755
/etc/init.d/single	0100755
