# var/lib/opkg/info/modutils-initscripts.list
/etc/init.d/modutils.sh	0100755
/etc/init.d	040755
/etc	040755
