# var/lib/opkg/info/openatv-bootlogo.list
/usr/share/bootlogo.mvi	0100644
/etc/init.d	040755
/usr/share/backdrop.mvi	0120777
/usr/share/restore.mvi	0100644
/usr	040755
/etc	040755
/usr/share/enigma2	040755
/etc/init.d/bootlogo	0100755
/usr/share	040755
/usr/share/enigma2/radio.mvi	0100644
