# var/lib/opkg/info/shellinabox.list
/etc/init.d	040755
/usr/bin/shellinaboxd	0100755
/usr/bin	040755
/etc/init.d/shellinabox	0100755
/usr	040755
/etc	040755
