# var/lib/opkg/info/sysvinit.list
/sbin/poweroff.sysvinit	0120777
/bin	040755
/usr/bin/last.sysvinit	0100755
/usr/bin/utmpdump.sysvinit	0100755
/etc/default/rcS	0100644
/usr/bin/mesg.sysvinit	0100755
/etc/init.d	040755
/etc/rc0.d	040755
/sbin/reboot.sysvinit	0120777
/etc/rc1.d	040755
/etc/init.d/rc	0100755
/etc/rc2.d	040755
/usr/bin	040755
/bin/mountpoint.sysvinit	0100755
/usr/bin/lastb.sysvinit	0120777
/etc/rc3.d	040755
/etc/rcS.d	040755
/sbin	040755
/sbin/bootlogd	0100755
/usr	040755
/etc	040755
/etc/rc4.d	040755
/etc/default/volatiles	040755
/etc/init.d/bootlogd	0100755
/sbin/runlevel.sysvinit	0100755
/etc/rc5.d	040755
/etc/init.d/stop-bootlogd	0120777
/sbin/telinit	0120777
/etc/default	040755
/etc/rc6.d	040755
/sbin/fstab-decode	0100755
/sbin/init.sysvinit	0100755
/usr/bin/wall.sysvinit	0100755
/etc/init.d/rcS	0100755
/sbin/shutdown.sysvinit	0104754
/sbin/halt.sysvinit	0104754
/usr/bin/readbootlog	0100755
/etc/default/volatiles/01_bootlogd	0100644
/etc/inittab.d	040755
