# var/lib/opkg/info/volatile-media.list
/etc/init.d	040755
/etc/init.d/volatile-media.sh	0100755
/etc/rcS.d	040755
/etc	040755
/etc/rcS.d/S02volatile-media.sh	0120777
